# cookbooks/redis/files/default/redis.service.conf
[program:redis]
user=docker
command=/usr/local/bin/redis-server /storage/redis.conf
stdout_events_enabled = true
stderr_events_enabled = true
