# cookbooks/redis/files/default/sentinel.service.conf
[program:sentinel]
user=docker
command=/usr/local/bin/redis-sentinel /storage/redis-sentinel.conf
