# redis.service.conf
[program:redis]
command=/usr/bin/redis-server /storage/redis/redis.conf
