# services/redis-sentinel.service.conf
[program:sentinel]
user=docker
command=/usr/bin/redis-sentinel /storage/redis/redis-sentinel.conf
