# services/redis.service.conf
[program:redis]
user=docker
command=/usr/bin/redis-server /storage/redis/redis.conf
stdout_events_enabled = true
stderr_events_enabled = true
